# crontab
* * * * * sleep 1; 
* * * * * sleep 2; /usr/cron.sh
* * * * * sleep 3; /usr/cron.sh
* * * * * sleep 4; /usr/cron.sh
* * * * * sleep 5; /usr/cron.sh
